# old/backup/web.cron
#!/bin/sh
curl --silent --compressed http://www.flatcap.org/cron.php
